# tp3/src/yt.service
[Unit]
Description=Dl des vidéos ytb.

[Service]
ExecStart=/srv/yt/yt-v2.sh
User=mat

[Install]
WantedBy=multi-user.target
